# pennypincher.service
[Unit]
Description=INS PennyPincher

[Service]
Type=oneshot
ExecStart=/usr/sbin/tuned-adm profile balanced

[Install]
WantedBy=multi-user.target
